# packaging/linux/deb/postinst
#!/bin/bash
# Post-installation script for simple-smtp-mailer
# Requires license acceptance

set -e

PROJECT_NAME="simple-smtp-mailer"
CONFIG_DIR="/etc/$simple-smtp-mailer"
LOG_DIR="/var/log/$simple-smtp-mailer"
SERVICE_USER="smtpmailerdev"

# Display license and require acceptance
if [ "$1" = "configure" ]; then
    echo "=========================================="
    echo "simple-smtp-mailer License Agreement"
    echo "=========================================="
    echo ""
    cat /usr/share/doc/$simple-smtp-mailer/copyright 2>/dev/null || \
        cat /usr/share/$simple-smtp-mailer/LICENSE.txt 2>/dev/null || \
        echo "Please review the license at: https://github.com/simpledaemons/$simple-smtp-mailer/blob/main/LICENSE"
    echo ""
    echo "By continuing, you agree to the terms of the license."
    echo ""
    read -p "Do you accept the license? [y/N]: " -n 1 -r
    echo
    if [[ ! $REPLY =~ ^[Yy]$ ]]; then
        echo "License not accepted. Installation cancelled."
        exit 1
    fi
    
    # Create service user
    if ! id "$SERVICE_USER" &>/dev/null; then
        useradd -r -s /usr/sbin/nologin -d /var/lib/$simple-smtp-mailer -c "$simple-smtp-mailer service user" "$SERVICE_USER"
    fi
    
    # Create directories
    mkdir -p "$CONFIG_DIR"
    mkdir -p "$LOG_DIR"
    mkdir -p "/var/lib/$simple-smtp-mailer"
    
    # Set permissions
    chown -R "$SERVICE_USER:$SERVICE_USER" "$CONFIG_DIR"
    chown -R "$SERVICE_USER:$SERVICE_USER" "$LOG_DIR"
    chown -R "$SERVICE_USER:$SERVICE_USER" "/var/lib/$simple-smtp-mailer"
    
    # Enable and start service
    systemctl daemon-reload
    systemctl enable "$simple-smtp-mailer" 2>/dev/null || true
    systemctl start "$simple-smtp-mailer" 2>/dev/null || true
fi

exit 0

